# meta-cross/recipes-support/drivapi-audio/files/drivapi-dashboard.service
[Unit]
Description=DrivaPi Dashboard (Qt)
After=weston.service pipewire-system.service pipewire-pulse-system.service
Wants=weston.service pipewire-system.service pipewire-pulse-system.service wireplumber-system.service

[Service]
Type=simple
Environment=WAYLAND_DISPLAY=wayland-0
Environment=QT_QPA_PLATFORM=wayland-egl
Environment=PULSE_SERVER=unix:/run/pipewire/pulse/native

# Wayland compositor must be ready before Qt can open a surface
ExecStartPre=/bin/sh -c 'test -S /run/wayland/wayland-0'
# Retry loop: wait up to 5s for PipeWire pulse socket
ExecStartPre=/bin/sh -c 'for i in $(seq 1 10); do [ -S /run/pipewire/pulse/native ] && exit 0; sleep 0.5; done; exit 1'

ExecStart=/usr/bin/myqtapp
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
